# scripts/bbb/push/b3lb-push.path
[Unit]
Description=triggering b3lb-push by filesystem changes

[Path]
PathChanged=/var/bigbluebutton/published/presentation
TriggerLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
