# cleaner/b3lb-cleaner.timer
[Unit]
Description=daily BBB meetings cleanup

[Timer]
OnCalendar=05:00

[Install]
WantedBy=multi-user.target
