# cleaner/b3lb-cleaner.service
[Unit]
Description=BBB meetings cleanup

[Service]
Type=oneshot
ExecStart=/opt/b3lb/scripts/cleaner.py
